# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:recorder]
command=bash -c "python -u main.py record --config ${NARCI_CONFIG:-configs/um_future_recorder.yaml} ${NARCI_SYMBOL:+--symbol $NARCI_SYMBOL}"
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=999
stopwaitsecs=10
stopsignal=INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:healthcheck]
command=python -u deploy/healthcheck.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
